# scripts/crontab
#!/bin/sh

* * * * * sh /var/www/html/scripts/shellFetch.sh >> /var/log/crond_fetch.log 2>&1
